# deploy/systemd/clove-relay.service
[Unit]
Description=Clove Relay Server
Documentation=https://github.com/clove-project/clove
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=clove
Group=clove

# Environment configuration
EnvironmentFile=-/etc/clove/relay.env

# Working directory
WorkingDirectory=/opt/clove/relay

# Main executable (runs both WebSocket and REST API)
ExecStart=/usr/bin/python3 /opt/clove/relay/relay_server.py --host 0.0.0.0 --port 8765

# Restart configuration
Restart=always
RestartSec=10
WatchdogSec=60

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=clove-relay

# Resource limits
MemoryMax=1G
CPUQuota=200%
TasksMax=1000

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes

# Allow write access to data directories
ReadWritePaths=/var/lib/clove

[Install]
WantedBy=multi-user.target
