# deploy/systemd/clove-kernel.service
[Unit]
Description=Clove Kernel - Microkernel OS for AI Agents
Documentation=https://github.com/clove-project/clove
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=clove
Group=clove

# Environment configuration
EnvironmentFile=-/etc/clove/config.env

# Working directory
WorkingDirectory=/opt/clove

# Main executable
ExecStart=/opt/clove/clove_kernel

# Restart configuration
Restart=always
RestartSec=10
WatchdogSec=60

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=clove-kernel

# Resource limits
MemoryMax=512M
CPUQuota=100%
TasksMax=100

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes

# Allow write access to specific directories
ReadWritePaths=/var/lib/clove
ReadWritePaths=/tmp/clove.sock

# Capabilities
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
