# deploy/systemd/clove-tunnel.service
[Unit]
Description=Clove Relay Tunnel Client
Documentation=https://github.com/clove-project/clove
After=clove-kernel.service
Requires=clove-kernel.service
BindsTo=clove-kernel.service

[Service]
Type=simple
User=clove
Group=clove

# Environment configuration
EnvironmentFile=-/etc/clove/config.env

# Working directory
WorkingDirectory=/opt/clove

# Main executable
ExecStart=/usr/bin/python3 /opt/clove/scripts/tunnel_client.py

# Restart configuration
Restart=always
RestartSec=5
WatchdogSec=120

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=clove-tunnel

# Resource limits
MemoryMax=128M
CPUQuota=25%

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes

# Allow read access to kernel socket
ReadWritePaths=/tmp/clove.sock

[Install]
WantedBy=multi-user.target
